# i3/audio-setup.service
[Unit]
Description=Runs script to fix samsung audio

[Service]
User=root
WorkingDirectory=/usr/local/sbin
ExecStart=necessary-verbs.sh
